# modules/scamper/files/upstart-scamper.conf
# Scamper deamon
start on runlevel [2345]
stop on shutdown

respawn
respawn limit 10 5

expect fork

exec /usr/bin/scamper -D -P 23456 -M $(hostname)
