# files/etc/update-motd.d/60-node-information
#!/usr/bin/env python
import os, re, socket, sys, simplejson, urllib, pycountry

API='https://ring.nlnog.net/api/1.0'
HOSTNAME=socket.gethostname()
BANNER='/etc/update-motd.d/ring.nlnog.net-logo'

def get_api_object(url):
    response = simplejson.load(urllib.urlopen(url))
    returncode = response['info']['success'];
    rescount   = response['info']['resultcount'];
    if (returncode != 1):
        sys.exit("unable to retrieve api object '"+url+"'")
    if (rescount < 1):
        return undef
    else:
        return response['results']

def get_admininfo(pid):
    node=socket.gethostbyaddr(HOSTNAME)[0]
    org = re.sub('..$','',node)
    return os.popen("getent passwd "+org+" | cut -d : -f 5").read().strip()

def get_countryname(country,state):
    cname=''
    try:
        c=pycountry.countries.get(alpha2=country)
        cname=c.name
    except KeyError:
        cname='Unknown'
    if state:
        try:
            s=pycountry.subdivisions.get(code=country+'-'+state)
            cname=s.name+', '+cname                
        except KeyError:
            cname='Unknown, '+cname
    return cname

def get_cymruinfo(asn):
    cymruinfo=os.popen("dig +short TXT AS"+asn+".asn.cymru.com").read()
    return re.sub('"','',cymruinfo,count=0)

def print_nodeinfo(asn,localadmin,location,cymruasinfo):

    # NLNOG RING Logo - can be disabled :)
    try:
        f = open(BANNER, 'r')
        print f.read(),
        f.close()
    except IOError:
        print "File" + sys.argv[1] + "does not exist."

    print ""
    print "  Welcome on "+HOSTNAME+", an NLNOG RING Node!"
    print "  System operated by "+localadmin
    if cymruasinfo:
	cymrulist=cymruasinfo.split("|")
	region=cymrulist[2].strip()
	asinfo=cymrulist[4].strip()
        print "  Location: "+location+" - AS"+asn+" ("+region+", "+asinfo+")"
    else:
        print "  Location: "+location+" - AS"+asn

    # Munin / Nagios / Smokeping
    print ""
    print "  Munin / Graphite:"
    print
    print "  https://ring.nlnog.net/munin/ring.nlnog.net/"+HOSTNAME
    print "  http://graphite.ring.nlnog.net/"

    print ""
    print "  For more information, please visit https://ring.nlnog.net/"

if __name__ == "__main__":    
    #url=API+'/nodes/'+HOSTNAME
    url=API+'/nodes/278'

    result=get_api_object(url)
    if not result:
        sys.exit("unable to retrieve node info")
    nodeinfo = result['nodes'][0]
    asn=str(nodeinfo['asn'])
    pid=nodeinfo['participant']
    country=nodeinfo['countrycode']
    state=nodeinfo['statecode']

    localadmin=get_admininfo(pid)
    location=get_countryname(country,state)
    cymruasinfo=get_cymruinfo(asn)

    print_nodeinfo(asn,localadmin,location,cymruasinfo)

